# scripts/com.football.prefetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.football.prefetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/CursorProject/football/scripts/prefetch.sh</string>
    <string>full</string>
  </array>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>/Users/CursorProject/football/logs/prefetch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/CursorProject/football/logs/prefetch.err</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
